# supervisord.conf
[supervisord]
nodaemon=true

[program:main]
command=python src/main.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/main.err.log
stdout_logfile=/var/log/main.out.log

[program:telegram_bot]
command=python src/telegram_bot.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/telegram_bot.err.log
stdout_logfile=/var/log/telegram_bot.out.log
